# source/mkimg-eswin/addons/80-workaround-lightdm-X-on-drm-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/libexec/kill-lightdm-X"
